# dist/profile/files/rhea/upstart-rhea.conf
# THIS FILE IS MANAGED BY PUPPET. CHANGES WILL BE OVERWRITTEN.
start on net-device-up and started docker
task
console log
exec /opt/vpm/puppet/bin/run --update
